# alpine-letsencrypt/root/etc/services.d/crond/run
#!/usr/bin/with-contenv sh

exec crond -f -d 8 -c /etc/crontabs -L /dev/stdout
